# pizza.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/pizza.sock

[Install]
WantedBy=sockets.target
